# growhaus-agent.service
[Unit]
Description=Grow Light Automation

[Service]
Type=simple
Restart=always
RestartSec=30
ExecStart=/home/pi/growhaus/bin/growhaus


